# lib/assets/crontab
MAIL='keburgett@gmail.com'
32 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 1 1
33 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 1 0
34 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 2 1
35 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 2 0
36 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 3 1
37 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 3 0
38 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 4 1
39 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 4 0
40 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 5 1
41 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 5 0
42 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 6 1
43 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 6 0
44 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 7 1
45 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 7 0
46 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 8 1
47 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 8 0
48 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 9 1
49 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 9 0
50 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 10 1
51 11 * * 4 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:2017 10 0
