# lib/assets/crontab
MAIL='keburgett@gmail.com'
05 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 1 1 5484 99
06 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 1 0 5484 99
07 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 2 1 4335 99
08 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 2 0 4335 99
09 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 3 1 1398 99
10 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 3 0 1398 99
11 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 4 1 7802 99
12 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 4 0 7802 99
13 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 5 1 16793 99
14 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 5 0 16793 99
15 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 6 1 14232 99
16 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 6 0 14232 99
17 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 7 1 30897 99
18 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 7 0 30897 99
19 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 8 1 25770 99
20 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 8 0 25770 99
21 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 9 1 5277 99
22 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 9 0 5277 99
23 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 10 1 7038 99
24 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 10 0 7038 99
