# lib/assets/crontab
MAIL='keburgett@gmail.com'
52 12 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  1 localhost:3000 1
53 12 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  0 localhost:3000 1
54 12 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  1 localhost:3000 2
55 12 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  0 localhost:3000 2
56 12 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  1 localhost:3000 3
57 12 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  0 localhost:3000 3
58 12 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  1 localhost:3000 4
59 12 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  0 localhost:3000 4
00 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  1 localhost:3000 5
01 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  0 localhost:3000 5
02 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  1 localhost:3000 6
03 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  1  0 localhost:3000 6
04 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  1 localhost:3000 7
05 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  2  0 localhost:3000 7
06 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  1 localhost:3000 8
07 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  3  0 localhost:3000 8
08 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  1 localhost:3000 9
09 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  4  0 localhost:3000 9
10 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  1 localhost:3000 10
11 13 * * 0 sh /home/kenb/development/HyperAQ/lib/tasks/valve_actuator.sh  5  0 localhost:3000 10
