# lib/assets/crontab
MAIL='keburgett@gmail.com'
42 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 1 1 7239 99
43 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 1 0 7239 99
44 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 2 1 14317 99
45 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 2 0 14317 99
46 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 3 1 5710 99
47 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 3 0 5710 99
48 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 4 1 21316 99
49 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 4 0 21316 99
50 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 5 1 24878 99
51 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 5 0 24878 99
52 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 6 1 19136 99
53 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 6 0 19136 99
54 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 7 1 26791 99
55 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 7 0 26791 99
56 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 8 1 27655 99
57 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 8 0 27655 99
58 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 9 1 1613 99
59 16 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 9 0 1613 99
00 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 10 1 27758 99
01 17 * * 2 sh /home/kenb/development/HyperAQ/lib/tasks/sprinkle_actuator.sh localhost:3000 10 0 27758 99
